# services/quantum/templates/etc/udev/rules.d/configure_br-ex.rules
SUBSYSTEM=="net", KERNEL=="br-ex", ACTION=="add", RUN+="/etc/udev/configure_br-ex"
